# deployment/systemd/lawim-backup.service
[Unit]
Description=LAWIM Backup and Disaster Recovery
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/lawim/current
EnvironmentFile=-/etc/lawim/backup.env
ExecStart=/bin/bash /opt/lawim/current/deployment/backup/backup.sh
TimeoutStartSec=3600
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
StateDirectory=lawim-backup
LogsDirectory=lawim-backup
RuntimeDirectory=lawim-backup
UMask=0077
PrivateTmp=true
NoNewPrivileges=true
