# deployment/systemd/lawim-backup.timer
[Unit]
Description=LAWIM twice-daily verified backup

[Timer]
OnCalendar=*-*-* 02:00:00 Africa/Douala
OnCalendar=*-*-* 14:30:00 Africa/Douala
Persistent=true
RandomizedDelaySec=5m
Unit=lawim-backup.service

[Install]
WantedBy=timers.target
